# deploy/hostinger/homostay_app_puma.service
[Unit]
Description=Homestay App Puma
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/homostay_app/server
EnvironmentFile=/etc/homostay_app.env
ExecStartPre=/usr/bin/mkdir -p /var/www/homostay_app/server/tmp/pids /var/www/homostay_app/server/tmp/sockets
ExecStart=/home/deploy/.rbenv/shims/bundle exec puma -C config/puma.rb
ExecReload=/bin/kill -USR1 $MAINPID
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=30
SyslogIdentifier=homostay_app_puma

[Install]
WantedBy=multi-user.target
